# deploy/bare-metal/godmode.service
[Unit]
Description=GodMode
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=godmode
Group=godmode
WorkingDirectory=/opt/godmode/current
EnvironmentFile=-/etc/godmode/godmode.env
Environment=PLATFORM_DATA_DIR=/var/lib/godmode
ExecStart=/opt/godmode/current/bin/godmode
Restart=on-failure
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/godmode

[Install]
WantedBy=multi-user.target
